# minichesstrain-stop.timer
[Unit]
Description=Mini Chess AI Training Stop Timer (6 AM America/Los_Angeles daily)
Documentation=file:///srv/MiniChessVovka/README_AI.md

[Timer]
# Stop at 6 AM (06:00) in America/Los_Angeles timezone
OnCalendar=*-*-* 06:00:00 America/Los_Angeles

# Run immediately on boot if we're past the stop time
Persistent=true

# Accuracy: trigger within 1 minute of scheduled time
AccuracySec=1min

[Install]
WantedBy=timers.target
