# minichesstrain.timer
[Unit]
Description=Mini Chess AI Training Timer (00:00 UTC daily start)
# НЕ используем Requires - пусть timer живет независимо от service

[Timer]
# Запускаем каждый день в 00:00 UTC
# Скрипт сам ждёт 02:00 UTC чтобы начать обучение и завершится в 10:00 UTC
OnCalendar=*-*-* 00:00:00

# Если пропустили (сервер был выключен), не запускать при следующей загрузке
Persistent=false

# Точность срабатывания
AccuracySec=1min

[Install]
WantedBy=timers.target
