# minichesstrain.service
[Unit]
Description=Mini Chess AI Self-Play Training
After=network.target
Wants=minichesstrain.timer

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/MiniChessVovka
ExecStart=/usr/bin/python3 /srv/MiniChessVovka/src/scheduled_self_play.py
Restart=on-failure
RestartPreventExitStatus=0
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=minichesstrain

# Resource Limits
CPUQuota=80%
MemoryMax=4G
MemorySwapMax=4G

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/MiniChessVovka

[Install]
WantedBy=multi-user.target
