# minichesstrain.timer
[Unit]
Description=Mini Chess AI Training Timer (0-8 UTC daily)
Requires=minichesstrain.service

[Timer]
# Проверяем каждые 30 минут
OnBootSec=5min
OnUnitActiveSec=30min

# Одноразовый запуск в 00:00 UTC (перезапускается через timer)
# Если хочешь в 0:00 ночи ежедневно:
# OnCalendar=daily
# OnCalendar=*-*-* 00:00:00 UTC

Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
